# playbooks/templates/lib/systemd/system/kubelet.service.d/50-crio.service.conf.j2
# {{ ansible_managed }}

[Unit]
After=crio.service
Wants=crio.service
