# playbooks/templates/lib/systemd/system/kubelet.service.d/50-etcd.service.conf.j2
# {{ ansible_managed }}

[Unit]
After=etcd.service
Wants=etcd.service
