# playbooks/templates/lib/systemd/system/kubelet.service.d/ceph-osd.conf.j2
# {{ ansible_managed }}

[Unit]
After=ceph-osd.target
Wants=ceph-osd.target
