# playbooks/templates/lib/systemd/system/crio.service.d/50-ceph.target.conf.j2
# {{ ansible_managed }}

[Unit]
After=ceph.target
Requires=ceph.target
