# playbooks/templates/lib/systemd/system/crio.service.d/ceph.conf.j2
# {{ ansible_managed }}

[Unit]
After=ceph.target
Wants=ceph.target
